# deursoos.service
[Unit]
Description=Software to open the door
 
[Service]
Type=simple
ExecStart=/opt/deursysteem/scan
Restart=always
