# deursoos.service
[Unit]
Description=Software to open the door
Requires=network-online.target
 
[Service]
Type=simple
ExecStart=/opt/deursysteem/scan
Restart=always
